# cronjobs/updater.cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
GEM_HOME=/usr/local/bundle

0 11 * * * echo "hello stackoverflow" > /proc/1/fd/1 2>/proc/1/fd/2
# * * * * * echo "hello stackoverflow" >> /var/log/cron.log 2>&1
1 11 * * * /app/script.sh > /proc/1/fd/1 2>/proc/1/fd/2
# * * * * * env > ~/cronenv

# * * * * * /app/ruby.rb > /proc/1/fd/1 2>/proc/1/fd/2

